# deployment/aclarknet.socket
# systemd socket unit for aclarknet Django/Wagtail application
# 
# This socket file implements socket activation for the Gunicorn WSGI server.
# Socket activation allows systemd to create and listen on the socket before
# the main service starts, enabling on-demand service activation and graceful
# restarts without dropping connections.
#
# The Unix socket serves as the communication channel between:
#   nginx (reverse proxy) → Unix socket → Gunicorn → Django application
#
# Related files:
#   - aclarknet.service: The main service that runs Gunicorn
#   - nginx-aclarknet.conf: nginx configuration that proxies to this socket
#
# To check status: sudo systemctl status aclarknet.socket
# To restart: sudo systemctl restart aclarknet.socket

[Unit]
Description=aclarknet gunicorn socket

[Socket]
# Create Unix socket for nginx to communicate with Gunicorn
ListenStream=/run/gunicorn/aclarknet.sock

# Create the /run/gunicorn directory with correct permissions
# This ensures the parent directory exists before the socket is created
RuntimeDirectory=gunicorn
RuntimeDirectoryMode=0755

# Set ownership to nginx user/group so nginx can read/write to the socket
SocketUser=nginx
SocketGroup=nginx

# Socket permissions: read/write for owner (nginx) and group (nginx)
SocketMode=0660

[Install]
# Start socket when entering multi-user mode (system boot)
WantedBy=sockets.target
